# roles/exporters/templates/filebeat.service.j2
[Unit]
Description=Filebeat-agent
After=network.target

[Service]
User=root
Group=root
ExecStart={{exporter_path}}/filebeat -e -c {{exporter_path}}/filebeat.yml

[Install]
WantedBy=multi-user.target