# roles/victoria-cluster/templates/node-exporter.service.j2
[Unit]
Description=Node-exporter
After=network.target

[Service]
User=root
Group=root
ExecStart={{victoria_path}}/node_exporter

[Install]
WantedBy=multi-user.target